# .config/systemd/user/copyq@.service
[Service]
Environment="DISPLAY=%I"
ExecStart=/usr/bin/copyq
TimeoutStopSec=0.1
Restart=on-failure
RestartSec=1s
