# .config/systemd/user/flameshot@.service
[Service]
Environment="DISPLAY=%I"
ExecStart=/usr/bin/flameshot
