# .config/systemd/user/kdeconnectd.service
[Service]
ExecStart=/usr/bin/kdeconnectd --replace

[Unit]
Requires=notifications.service
