# .config/systemd/user/kdeconnectd@.service
[Service]
Environment="DISPLAY=%I"
ExecStart=/usr/lib/kdeconnectd --replace
