# .config/systemd/user/kdeconnectd.service
[Service]
ExecStart=/usr/lib/kdeconnectd --replace
