# .config/systemd/user/qmenu_hud.service
[Unit]
Description=qmenu HUD daemon

[Service]
ExecStart=/usr/bin/qmenu_registrar
