# .config/systemd/user/qmenu_hud@.service
[Unit]
Description=qmenu HUD daemon

[Service]
Environment="DISPLAY=%I"
ExecStart=/usr/bin/qmenu_registrar
