# .config/systemd/user/solaar@.service
[Service]
Environment="DISPLAY=%I"
ExecStart=/usr/bin/solaar -w hide
