# .config/systemd/user/sxhkd@.service
[Service]
Environment="DISPLAY=%I"
ExecStart=/usr/bin/sxhkd
ExecReload=/usr/bin/pkill -USR1 -x sxhkd
