# .config/systemd/user/vonal.service
[Service]
ExecStart=/home/dbiro/Projects/All/Programming/vonal/vonal/target/release/vonal
